# root/etc/services.d/app/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc command
  > /dev/stdout \
  2 > /dev/stderr

echo 'job'